# deploy/watchdog/watchdog.service
[Unit]
Description=Portfolio watchdog — probe prod and auto-recover
Documentation=https://github.com/Elkhair/portfolio-project/tree/master/deploy/watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/watchdog/check.sh
# A recovery tick can spend ~3 min waiting for a rebooted VM's SSH to come
# up + run `docker compose up -d`. Default 90s would kill the tick mid-wait.
TimeoutStartSec=300
# Lock down a bit — script only needs network + /var/lib/watchdog
ProtectSystem=strict
ReadWritePaths=/var/lib/watchdog /tmp
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
